# snake.conf
[program:snake]
command=python3 snake.py
directory=/home/pi/retro-pysnake
autostart=true
autorestart=true
environment=DISPLAY=":0.0"
